# embark.service
[Unit]
Description=EMBArk supervisor daemon

[Service]
Type=simple
ExecStart=/root/amos/supervisor.sh

[Install]
WantedBy=multi-user.target