# embark.service
[Unit]
Description=EMBArk supervisor daemon

[Service]
Type=simple
ExecStart=<BASE_DIR>/supervisor.sh

[Install]
WantedBy=multi-user.target