# embark.service
[Unit]
Description=EMBArk supervisor daemon

[Service]
Type=simple
ExecStart=/home/benedikt/git-repos/embark/supervisor.sh

[Install]
WantedBy=multi-user.target