# embark.service
[Unit]
Description=EMBArk supervisor daemon

[Service]
Type=simple
ExecStart={$EMBARK_ROOT_DIR}/supervisor.sh

[Install]
WantedBy=multi-user.target